# arch-packages/telegram-interaction/usr_lib_systemd_system_interaction.service
[Unit]
Description=Interaction Service

[Service]
Type=dbus
User=interaction
ExecStart=/usr/bin/interactiond
BusName=com.github.radium226.Interaction
