# arch-packages/maintenance/usr_lib_systemd_system_maintenance.service
[Unit]
Description=Maintenance Service

[Service]
Type=dbus
User=maintenance
ExecStart=/usr/bin/maintenanced --system
BusName=com.github.radium226.Maintenance
