# arch-packages/system-maintenance/usr_lib_systemd_system_system-maintenance@.service
[Unit]
Description=System Maintenance

[Service]
Type=oneshot
User=system-maintenance
ExecStart=/usr/bin/python -m "system.maintenance" "%i"
