# arch-packages/whistleblower/usr_lib_systemd_system_whistleblower.service
[Unit]
Description=Whistleblower

[Service]
Type=dbus
User=root
ExecStart=/usr/bin/whistleblowerd
BusName=com.github.radium226.Whistleblower
